# boodle.service
[Unit]
Description=Service to run an instance of boodle

[Service]
WorkingDirectory=/home/manuel/boodle
ExecStart=deploy.sh
